# environments/kolla/files/latest-full-mariadb-backup.timer
[Unit]
Description=Create full mariadb backup every 6 hours

[Timer]
OnCalendar=*-*-* 05/6:00:00
Persistent=true

[Install]
WantedBy=timers.target
