# playbooks/managerless/templates/sshuttle.service.j2
[Unit]
Description=OSISM testbed sshuttle service
After=network.target

[Service]
Type=simple
RemainAfterExit=true
ExecStart=/usr/local/bin/sshuttle.sh

[Install]
WantedBy=multi-user.target
